# auth.microtv.ar/deploy/auth-microtv.service
[Unit]
Description=auth.microtv.ar backend (uvicorn)
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=authmicrotv
Group=authmicrotv
WorkingDirectory=/opt/microtv/auth.microtv.ar/backend
EnvironmentFile=/opt/microtv/auth.microtv.ar/backend/.env
ExecStart=/opt/microtv/auth.microtv.ar/backend/venv/bin/uvicorn src.main:app \
    --host 0.0.0.0 \
    --port 8001 \
    --workers 4 \
    --log-level info \
    --access-log \
    --proxy-headers \
    --forwarded-allow-ips="0.0.0.0"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=auth-microtv

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/microtv/auth.microtv.ar/backend
PrivateTmp=true

[Install]
WantedBy=multi-user.target
